# deploy/bar-pos-backup.service
[Unit]
Description=bar-pos database backup (pg_dump with retention)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/bar-pos-backup.sh
